# rpi4/h10.service
[Unit]
Description=H10 — Polar H10 Heart-Rate BLE Relay Worker (bpm / RR intervals)
After=network.target bluetooth.target

[Service]
Type=simple
User=${SERVICE_USER}
WorkingDirectory=${WORKING_DIR}
ExecStart=${PYTHON_BIN} ${WORKING_DIR}/h10.py
Restart=on-failure
RestartSec=5
TimeoutStopSec=10

# Give the Bluetooth stack time to be ready after boot
ExecStartPre=/bin/sleep 5

# BLE access — ensure ${SERVICE_USER} is in the bluetooth group:
#   sudo usermod -aG bluetooth ${SERVICE_USER}
SupplementaryGroups=bluetooth

[Install]
WantedBy=multi-user.target
