# rpi4/pulse.service
[Unit]
Description=Pi-Pulse — System Metrics API (CPU / Memory / Temperature)
After=network.target

[Service]
Type=simple
User=${SERVICE_USER}
WorkingDirectory=${WORKING_DIR}
ExecStart=${UVICORN_BIN} pulse:app --host 0.0.0.0 --port 8001
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
