# rpi4/pulse.service
[Unit]
Description=Pi-Pulse — System Metrics Relay Worker (CPU / Memory / Temperature)
After=network.target

[Service]
Type=simple
User=${SERVICE_USER}
WorkingDirectory=${WORKING_DIR}
ExecStart=${PYTHON_BIN} ${WORKING_DIR}/pulse.py
Restart=on-failure
RestartSec=5
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
