# rpi4/sen66.service
[Unit]
Description=Sen66 ŌĆö Air Quality Relay Worker (PM / CO2 / VOC / NOx / RH / T)
After=network.target

[Service]
Type=simple
User=${SERVICE_USER}
WorkingDirectory=${WORKING_DIR}
ExecStart=${PYTHON_BIN} ${WORKING_DIR}/sen66.py
Restart=on-failure
RestartSec=5
TimeoutStopSec=10

# Give the I2C bus time to be ready after boot
ExecStartPre=/bin/sleep 3

# I2C bus access ŌĆö ensure ${SERVICE_USER} is in the i2c group:
#   sudo usermod -aG i2c ${SERVICE_USER}
SupplementaryGroups=i2c

[Install]
WantedBy=multi-user.target
